# docker/obsidian-kasm/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:obsidian]
command=/init
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/obsidian.log
stderr_logfile=/var/log/supervisor/obsidian-error.log

[program:sync-daemon]
command=node /opt/sync-daemon/dist/index.js
directory=/opt/sync-daemon
autostart=true
autorestart=true
priority=20
environment=NODE_ENV="production",SYNC_SERVER_URL="%(ENV_SYNC_SERVER_URL)s",VAULT_PATH="%(ENV_VAULT_PATH)s"
stdout_logfile=/var/log/supervisor/sync-daemon.log
stderr_logfile=/var/log/supervisor/sync-daemon-error.log
